# mouse-hardware-fix@.service
[Unit]
Description=Mouse Hardware Fix dynamically attached to %i
# Tie the lifecycle of this service directly to the hardware being plugged in
BindsTo=dev-input-%i.device
After=dev-input-%i.device

[Service]
Type=simple
# %i will be automatically replaced by the event node (e.g., event4)
ExecStart=/usr/bin/mouse-hardware-fix --device /dev/input/%I
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
